# upstart-redis-statsd.conf
description "Reporting redis metrics to local statsd"

start on runlevel [2345]

stop on runlevel [016]

exec /usr/local/sbin/redis-statsd.py --period 10

respawn
